# scripts/templates/systemd/eval-drift-daemon.service
[Unit]
Description=31C eval-drift daemon (nightly Langfuse trace replay)
After=network.target

[Service]
Type=simple
WorkingDirectory={{WORKSPACE}}
ExecStart={{PYTHON}} scripts/eval-drift-daemon.py daemon
Restart=on-failure
RestartSec=30
# See bridge-daemon.service for rationale: journald avoids double-writing the
# daemon's log path on a /mnt/c 9P mount, which can freeze the process in
# uninterruptible sleep.
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
